# usr/src/cmd/cmd-inet/usr.lib/wpa_supplicant/wpa_supplicant.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
	Copyright 2007 Sun Microsystems, Inc.  All rights reserved.
	Use is subject to license terms.
	Copyright (c) 2012, Enrico Papi. All rights reserved.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.
-->

<service_bundle type='manifest' name='SUNWsupr:wpa_supplicant'>

<service
	name='network/wpa_supplicant'
	type='service'
	version='1'>

	<single_instance/>

	<dependency
	    name='filesystem-minimal'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
	    <service_fmri value='svc:/system/filesystem/minimal' />
	</dependency>

	<dependency
	    name='cryptosvc'
	    grouping='optional_all'
	    restart_on='none'
	    type='service'>
	    <service_fmri value='svc:/system/cryptosvc' />
	</dependency>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/svc-wpa_supplicant'
		timeout_seconds='60'>
			<method_context>
				<method_credential
				user='root'
				group='netadm'
				limit_privileges=':default'
				privileges='basic,sys_net_config,net_rawaccess'
				/>
			</method_context>
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec=':kill'
		timeout_seconds='60'>
			<method_context>
				<method_credential
				user='root'
				group='netadm'
				limit_privileges=':default'
				privileges='basic,sys_net_config,net_rawaccess'
				/>
			</method_context>
	</exec_method>

  	<!-- reparse wpa_s configuration when using nwam -->
	<exec_method
		type='method'
		name='refresh'
		exec=':kill -HUP'
		timeout_seconds='60' />

        <property_group name='general' type='framework'>
                <!-- to start stop wpa_supplicant service
                	(as of now only dladm is allowed) -->
                <propval name='action_authorization' type='astring'
                        value='solaris.smf.manage.wpa' />
                <propval name='value_authorization' type='astring'
			value='solaris.smf.value.wpa' />
        </property_group>

	<property_group name='config' type='application' >
		<!-- default property settings for wpa_supplicant -->
		<propval
		    name='debug_enabled'
		    type='boolean'
		    value='false' />

		<propval
		    name='verbose_debug'
		    type='boolean'
		    value='false' />

		<propval
		    name='keys_debug'
		    type='boolean'
		    value='false' />

		<propval
		    name='timestamps_debug'
		    type='boolean'
		    value='false' />

		<!-- to change properties -->
		<propval
		    name='value_authorization'
		    type='astring'
		    value='solaris.smf.value.wpa' />
	</property_group>

       	<instance name="default" enabled="true" />

	<stability value='External' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			Wireless WPA Supplicant Daemon Version 1.1
			</loctext>
		</common_name>
		<documentation>
			<doc_link
			    name='WPA Supplicant README'
				uri='http://w1.fi/gitweb/gitweb.cgi?p=hostap-1.git;a=blob_plain;f=wpa_supplicant/README;hb=HEAD'
				/>
		</documentation>
	</template>
</service>

</service_bundle>
